# packages/systemd/supernode.service
[Unit]
Description=N2N network edge instance
Documentation=man:edge(8)
After=network.target

[Service]
DynamicUser=Yes
ExecStart=supernode -t /run/supernode/mgmt -f -l $PORT $PARAMETER
EnvironmentFile=/etc/n2n/supernode.conf
RuntimeDirectory=supernode

[Install]
WantedBy=multi-user.target
